# detections/sigma/windows/persistence/runkeys_regexe_fallback.yml
title: Run/RunOnce Persistence via reg.exe (Fallback)
id: 2b9d6f13-8c3e-4bc6-9b8e-0a7f6e5d4c3b
status: experimental
description: Fallback detection for Run/RunOnce changes via reg.exe arguments.
author: Khaireddine
date: 2025/08/21
logsource:
  product: windows
  category: process_creation
  service: sysmon
detection:
  sel_img:
    Image|endswith: '\reg.exe'
  sel_cmd:
    CommandLine|contains:
      - '\CurrentVersion\Run'
      - '\CurrentVersion\RunOnce'
  condition: sel_img and sel_cmd
fields: [UtcTime, Image, CommandLine, ParentImage, User]
falsepositives:
  - Legitimate installers/IT scripts
level: medium
tags:
  - attack.persistence
  - attack.t1547
  - attack.t1547.001
