# detections/sigma/windows/persistence/schtasks_create_sysmon.yml
title: Scheduled Task Creation via schtasks.exe (Sysmon)
id: f7a9d3b1-6c4e-4f2a-b8e3-9c2d5f7a4b3e
status: experimental
description: Detects creation of scheduled tasks using schtasks.exe with the /create parameter.
author: Khaireddine
date: 2025/08/18
logsource:
  product: windows
  category: process_creation
  service: sysmon
detection:
  sel_img:
    Image|endswith: '\schtasks.exe'
  sel_cmd:
    CommandLine|contains: ' /create '
  condition: sel_img and sel_cmd
fields: [ UtcTime, Image, CommandLine, ParentImage, ParentCommandLine, User ]
falsepositives: [ Legitimate administrative scheduled task creation ]
level: high
tags: [ attack.persistence, attack.t1053, attack.t1053.005 ]
