# docs/mitre_layer_apts.json
{
  "name": "APT29 & Lazarus – Detection Coverage (ATT&CK v17, Windows)",
  "versions": {
    "attack": "17",
    "navigator": "4.9.1",
    "layer": "4.5"
  },
  "domain": "enterprise-attack",
  "description": "Coverage snapshot of 11 implemented detections mapped to ATT&CK v17. Color key: 🟩 Covered by working KQL/Sigma. Comments include severity and brief note.",
  "filters": {
    "platforms": ["Windows"]
  },
  "sorting": 0,
  "layout": {
    "layout": "flat",
    "showID": true,
    "showName": true,
    "showAggregateScores": false,
    "countUnscored": false
  },
  "hideDisabled": false,
  "legendItems": [
    { "label": "Covered by your KQL/Sigma", "color": "#2ECC71" }
  ],
  "gradient": {
    "colors": ["#ffffff", "#2ECC71"],
    "minValue": 0,
    "maxValue": 1
  },
  "metadata": [
    { "name": "Scope", "value": "Windows endpoint (Sysmon/Elastic Winlogbeat)" },
    { "name": "Owner", "value": "Your detection engineering project" }
  ],
  "links": [],
  "techniques": [
    {
      "techniqueID": "T1047",
      "tactic": "execution",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Windows Management Instrumentation (WMI) execution — Severity: Medium — Covered via process/command-line patterns."
    },
    {
      "techniqueID": "T1059.001",
      "tactic": "execution",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Command and Scripting Interpreter: PowerShell — Severity: High — Covers suspicious flags/args and LOLBIN handoffs."
    },
    {
      "techniqueID": "T1053.005",
      "tactic": "persistence",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Scheduled Task/Job: Scheduled Task — Severity: Medium — Creation/invocation patterns monitored."
    },
    {
      "techniqueID": "T1105",
      "tactic": "command-and-control",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Ingress Tool Transfer — Severity: High — Download/transfer behaviors via native tools covered."
    },
    {
      "techniqueID": "T1003.002",
      "tactic": "credential-access",
      "score": 1,
      "color": "#2ECC71",
      "comment": "OS Credential Dumping: Security Account Manager (SAM) — Severity: Critical — Access/dumping patterns covered."
    },
    {
      "techniqueID": "T1003.004",
      "tactic": "credential-access",
      "score": 1,
      "color": "#2ECC71",
      "comment": "OS Credential Dumping: LSA Secrets — Severity: Critical — Sensitive registry/LSA access behaviors covered."
    },
    {
      "techniqueID": "T1573",
      "tactic": "command-and-control",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Encrypted Channel — Severity: Medium — Indicators of encrypted/obfuscated C2 channel usage monitored."
    },
    {
      "techniqueID": "T1055.001",
      "tactic": "defense-evasion",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Process Injection: Dynamic-link Library Injection — Severity: High — Sysmon EID 7/8 (image load/remote thread) logic implemented."
    },
    {
      "techniqueID": "T1218.005",
      "tactic": "defense-evasion",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Signed Binary Proxy Execution: Mshta — Severity: Medium — HTA/URL/parent-child patterns covered."
    },
    {
      "techniqueID": "T1218.010",
      "tactic": "defense-evasion",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Signed Binary Proxy Execution: Regsvr32 — Severity: Medium — .sct/scrobj and suspicious /i usage covered."
    },
    {
      "techniqueID": "T1218.011",
      "tactic": "defense-evasion",
      "score": 1,
      "color": "#2ECC71",
      "comment": "Signed Binary Proxy Execution: Rundll32 — Severity: Medium — pcwutl.dll, url.dll handlers, Control_RunDLL, screensaver abuse covered."
    }
  ]
}
