# SigmaRules/persistence/scheduled_task.yml
title: Suspicious Scheduled Task Creation via schtasks.exe
id: f7a9d3b1-6c4e-4f2a-b8e3-9c2d5f7a4b3e
status: stable
description: Detects creation of scheduled tasks using schtasks.exe with /create parameter, considering case variations, and suspicious parent processes.
author: Khaireddine
date: 2025/09/05
logsource:
  product: windows
  category: process_creation
  service: sysmon
detection:
  sel_img:
    Image|endswith: '\schtasks.exe'
  sel_create:
    CommandLine|contains:
      - '/create'
      - '/Create'
  sel_tn:
    CommandLine|contains:
      - '/tn'
      - '/TN'
  sel_tr:
    CommandLine|contains:
      - '/tr'
      - '/TR'
  sel_parent:
    ParentImage|endswith:
      - '\cmd.exe'
      - '\powershell.exe'
  condition: sel_img and sel_create and sel_tn and sel_tr and sel_parent
fields:
  - UtcTime
  - Image
  - CommandLine
  - ParentImage
  - ParentCommandLine
  - User
falsepositives:
  - Legitimate administrative task creation
  - Enterprise management or installer tools
level: high
tags:
  - attack.persistence
  - attack.t1053
  - attack.t1053.005